# ops/services/chyren-brain.service
[Unit]
Description=Chyren Brain (Canonical Sovereign API)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mega
Group=mega
WorkingDirectory=/home/mega/Chyren
EnvironmentFile=-/home/mega/Chyren/state/vault/one-true.env
ExecStart=/home/mega/Chyren/core/chyren-os/kernel/target/release/chyren_api
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=5
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
