# ops/services/chyren-summary.service
[Unit]
Description=chyren Global Summary
After=network.target redis-server.service chyren-gateway.service chyren-bridge.service chyren-brain.service

[Service]
Type=oneshot
Environment=SUMMARY_OUT=/home/mega/Chyren/state/last_summary.md
ExecStart=/usr/bin/env bash -c '/home/mega/NEXUS/repos/chyren-Trinity/scripts/agent-bridge.sh summarize'

[Install]
WantedBy=multi-user.target
