# ops/services/chyren-revenue-router.service
[Unit]
Description=OMEGA Revenue Router — closed-loop income orchestration engine
After=network.target omega-gateway.service omega-brain.service

[Service]
Type=simple
User=mega
WorkingDirectory=/home/mega/OmegA-Sovereign/gAIng-brAin
EnvironmentFile=-/home/mega/OmegA-Sovereign/gAIng-brAin/.env
Environment=REVENUE_ROUTER_PORT=9091
ExecStart=/home/mega/.nvm/versions/node/v22.22.0/bin/node /home/mega/OmegA-Sovereign/gAIng-brAin/revenue-router.js
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
