# ops/services/chyren-revenue-router.service
[Unit]
Description=Chyren Revenue Router (Closed-Loop Orchestration)
After=chyren-brain.service

[Service]
Type=simple
User=mega
WorkingDirectory=/home/mega/Chyren
EnvironmentFile=-/home/mega/Chyren/state/vault/one-true.env
ExecStart=/home/mega/Chyren/chyren.legacy flex
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
