# infrastructure/legacy_scripts/services/chyren-telegram.service
[Unit]
Description=Chyren Telegram Gateway
After=network.target

[Service]
Type=simple
User=mega
WorkingDirectory=/home/mega/Chyren
EnvironmentFile=-/home/mega/Chyren/state/vault/one-true.env
ExecStart=/home/mega/Chyren/chyren telegram
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
