# infrastructure/legacy_scripts/services/chyren-gateway.service
[Unit]
Description=Chyren Gateway (Web Interface)
After=network.target chyren-brain.service
Requires=chyren-brain.service

[Service]
Type=simple
User=mega
WorkingDirectory=/home/mega/Chyren
EnvironmentFile=-/home/mega/Chyren/state/vault/one-true.env
ExecStart=/home/mega/Chyren/chyren.legacy web
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
