# infrastructure/legacy_scripts/services/chyren-workflows.service
[Unit]
Description=Chyren Workflows (Autonomous Logic)
After=chyren-brain.service

[Service]
Type=simple
User=mega
WorkingDirectory=/home/mega/Chyren
EnvironmentFile=-/home/mega/Chyren/state/vault/one-true.env
ExecStart=/home/mega/Chyren/chyren.legacy action
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
