# infrastructure/legacy_scripts/services/chyren-heartbeat.service
[Unit]
Description=chyren Heartbeat — feeds ERGON.md and task-queue into gateway memory
After=chyren-gateway.service

[Service]
Type=oneshot
User=mega
EnvironmentFile=-/home/mega/Chyren/state/vault/one-true.env
ExecStart=/bin/bash /home/mega/bin/chyren-heartbeat.sh
StandardOutput=journal
StandardError=journal
