# ops/chyren-dream.timer
[Unit]
Description=Run Chyren Dream Cycle Nightly

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
Unit=chyren-dream.service

[Install]
WantedBy=timers.target
